# personal/launchagents/com.vaf.morning-briefing.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.vaf.morning-briefing</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>/Users/mcmehmios/repos/v-agentforce-architecture/personal/scripts/morning-briefing.py</string>
    </array>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>7</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    <key>StandardOutPath</key>
    <string>/Users/mcmehmios/repos/v-agentforce-architecture/logs/morning-briefing.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/mcmehmios/repos/v-agentforce-architecture/logs/morning-briefing.error.log</string>
    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
