# vaishali-agent-force/macos/com.vaishali.agentforce.dashboard.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.vaishali.agentforce.dashboard</string>

    <key>ProgramArguments</key>
    <array>
        <string>__VAF_PYTHON__</string>
        <string>-m</string>
        <string>uvicorn</string>
        <string>vaishali.dashboard.api:app</string>
        <string>--host</string>
        <string>0.0.0.0</string>
        <string>--port</string>
        <string>8077</string>
    </array>

    <key>WorkingDirectory</key>
    <string>__VAF_PROJECT__</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>VAF_BASE_DIR</key>
        <string>__VAF_PROJECT__</string>
        <key>PYTHONPATH</key>
        <string>__VAF_PROJECT__/src</string>
    </dict>

    <!-- Start on login, keep running -->
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>__VAF_PROJECT__/logs/dashboard.log</string>
    <key>StandardErrorPath</key>
    <string>__VAF_PROJECT__/logs/dashboard.err</string>

    <key>Nice</key>
    <integer>5</integer>
</dict>
</plist>
